# deployment/georisk-frontend.service
[Unit]
Description=GeoRisk Frontend Server
After=network.target

[Service]
Type=simple
User=bitnami
Group=bitnami
WorkingDirectory=/home/bitnami/my-app
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/bitnami/.nvm/versions/node/v20.11.0/bin"
Environment="NODE_ENV=production"
ExecStart=/usr/bin/npx serve -s dist -l 5173
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
